# config/systemd/user/xscreensaver.service
[unit]
Description=Xscreensaver
[Service]
ExecStart=/usr/bin/xscreensaver
[Install]
WantedyBy=default.target
